# systemd/auto-repair.service
[Unit]
Description=AwesomeAgent Auto-Repair Service

[Service]
Type=oneshot
ExecStart=/opt/awesomeagent/scripts/auto_repair.sh
