# systemd/auto-update.service
[Unit]
Description=AwesomeAgent Auto-Update Service

[Service]
Type=oneshot
ExecStart=/opt/awesomeagent/scripts/auto_update.sh
